# templates/public-inbox-config-watcher.path.erb
[Path]
Unit=public-inbox-config-watcher.service
PathChanged=<%= @config %>

[Install]
WantedBy=multi-user.target
