# templates/public-inbox-watch.service.erb
# MANAGED BY PUPPET

[Unit]
Description = public-inbox Maildir watch
After = spamassassin.service

[Service]
Environment = PI_CONFIG=<%= @config %> \
              PATH=/usr/local/bin:/usr/bin:/bin
ExecStart = /usr/local/bin/public-inbox-watch

StandardOutput = syslog
StandardError = syslog
ExecReload = /bin/kill -HUP $MAINPID
# this user must have read access to Maildirs it watches
User = <%= @var_dir_owner %>
KillMode = process

[Install]
WantedBy = multi-user.target
